# alertmanager.ini
[program:alertmanager]
priority = 99
user = pmm
command =
    /usr/local/bin/alertmanager
        -config.file=/etc/alertmanager.yml
        -log.level=info
        -mesh.listen-address=''
        -storage.path=/opt/alertmanager/data/
        -web.external-url=http://localhost:9093/alertmanager/
        -web.listen-address=:9093
stdout_logfile = /var/log/alertmanager.log
stderr_logfile = /var/log/alertmanager.log
autorestart = true
stopwaitsecs = 300
